# Config/supervisor/server_133.conf
[group:Servers]
programs=AutoRestart,GetInfoFromLBK,GetInfoFromCMC,GetInfoFromUNIv2,GetInfoFromUNIv3,GetInfoFromPANCAKEv2,GetInfoFromPANCAKEv3,LBKWssData,SendMail,DailyReport,DailyReportMail,MonitorSTG,TradeCacheSTG

[program:AutoRestart]
directory=/home/ec2-user/MMSys-ws/Servers
command=/home/ec2-user/anaconda3/envs/lbk-tornado/bin/python auto_restart.py
autostart=true
autorestart=unexpected
startretries=5
startsecs=1
priority=999
stderr_logfile=/dev/null
stdout_logfile=/dev/null
stopasgroup=true
killasgroup=true
user = root
redirect_stderr = true
asgfjfghfc_maxbytes = 20M
stdout_logfile_backups = 20

[program:GetInfoFromLBK]
directory=/home/ec2-user/MMSys-ws/Servers
command=/home/ec2-user/anaconda3/envs/lbk-tornado/bin/python get_info_from_lbk.py
autostart=true
autorestart=unexpected
startretries=5
startsecs=1
priority=999
stderr_logfile=/dev/null
stdout_logfile=/dev/null
stopasgroup=true
killasgroup=true
user = root
redirect_stderr = true
asgfjfghfc_maxbytes = 20M
stdout_logfile_backups = 20

[program:GetInfoFromCMC]
directory=/home/ec2-user/MMSys-ws/Servers
command=/home/ec2-user/anaconda3/envs/lbk-tornado/bin/python get_info_from_cmc.py
autostart=true
autorestart=unexpected
startretries=5
startsecs=1
priority=999
stderr_logfile=/dev/null
stdout_logfile=/dev/null
stopasgroup=true
killasgroup=true
user = root
redirect_stderr = true
asgfjfghfc_maxbytes = 20M
stdout_logfile_backups = 20

[program:GetInfoFromUNIv2]
directory=/home/ec2-user/MMSys-ws/Servers
command=/home/ec2-user/anaconda3/envs/lbk-tornado/bin/python get_info_from_uni_v2.py
autostart=true
autorestart=unexpected
startretries=5
startsecs=1
priority=999
stderr_logfile=/dev/null
stdout_logfile=/dev/null
stopasgroup=true
killasgroup=true
user = root
redirect_stderr = true
asgfjfghfc_maxbytes = 20M
stdout_logfile_backups = 20

[program:GetInfoFromUNIv3]
directory=/home/ec2-user/MMSys-ws/Servers
command=/home/ec2-user/anaconda3/envs/lbk-tornado/bin/python get_info_from_uni_v3.py
autostart=true
autorestart=unexpected
startretries=5
startsecs=1
priority=999
stderr_logfile=/dev/null
stdout_logfile=/dev/null
stopasgroup=true
killasgroup=true
user = root
redirect_stderr = true
asgfjfghfc_maxbytes = 20M
stdout_logfile_backups = 20

[program:GetInfoFromPANCAKEv2]
directory=/home/ec2-user/MMSys-ws/Servers
command=/home/ec2-user/anaconda3/envs/lbk-tornado/bin/python get_info_from_pancake_v2.py
autostart=true
autorestart=unexpected
startretries=5
startsecs=1
priority=999
stderr_logfile=/dev/null
stdout_logfile=/dev/null
stopasgroup=true
killasgroup=true
user = root
redirect_stderr = true
asgfjfghfc_maxbytes = 20M
stdout_logfile_backups = 20

[program:GetInfoFromPANCAKEv3]
directory=/home/ec2-user/MMSys-ws/Servers
command=/home/ec2-user/anaconda3/envs/lbk-tornado/bin/python get_info_from_pancake_v3.py
autostart=true
autorestart=unexpected
startretries=5
startsecs=1
priority=999
stderr_logfile=/dev/null
stdout_logfile=/dev/null
stopasgroup=true
killasgroup=true
user = root
redirect_stderr = true
asgfjfghfc_maxbytes = 20M
stdout_logfile_backups = 20

[program:LBKWssData]
directory=/home/ec2-user/MMSys-ws/Servers
command=/home/ec2-user/anaconda3/envs/lbk-tornado/bin/python lbk_wss_data.py
autostart=true
autorestart=unexpected
startretries=5
startsecs=1
priority=999
stderr_logfile=/dev/null
stdout_logfile=/dev/null
stopasgroup=true
killasgroup=true
user = root
redirect_stderr = true
asgfjfghfc_maxbytes = 20M
stdout_logfile_backups = 20

[program:SendMail]
directory=/home/ec2-user/MMSys-ws/Servers
command=/home/ec2-user/anaconda3/envs/lbk-tornado/bin/python sendMail.py
autostart=true
autorestart=unexpected
startretries=5
startsecs=1
priority=999
stderr_logfile=/dev/null
stdout_logfile=/dev/null
stopasgroup=true
killasgroup=true
user = root
redirect_stderr = true
asgfjfghfc_maxbytes = 20M
stdout_logfile_backups = 20

[program:DailyReport]
directory=/home/ec2-user/MMSys-ws/Servers
command=/home/ec2-user/anaconda3/envs/lbk-tornado/bin/python dailyReport.py
autostart=true
autorestart=unexpected
startretries=5
startsecs=1
priority=999
stderr_logfile=/dev/null
stdout_logfile=/dev/null
stopasgroup=true
killasgroup=true
user = root
redirect_stderr = true
asgfjfghfc_maxbytes = 20M
stdout_logfile_backups = 20

[program:DailyReportMail]
directory=/home/ec2-user/MMSys-ws/Servers
command=/home/ec2-user/anaconda3/envs/lbk-tornado/bin/python dailyReport_mail.py
autostart=true
autorestart=unexpected
startretries=5
startsecs=1
priority=999
stderr_logfile=/dev/null
stdout_logfile=/dev/null
stopasgroup=true
killasgroup=true
user = root
redirect_stderr = true
asgfjfghfc_maxbytes = 20M
stdout_logfile_backups = 20

[program:MonitorSTG]
directory=/home/ec2-user/MMSys-ws/Servers
command=/home/ec2-user/anaconda3/envs/lbk-tornado/bin/python monitor_stg.py
autostart=true
autorestart=unexpected
startretries=5
startsecs=1
priority=999
stderr_logfile=/dev/null
stdout_logfile=/dev/null
stopasgroup=true
killasgroup=true
user = root
redirect_stderr = true
asgfjfghfc_maxbytes = 20M
stdout_logfile_backups = 20

[program:TradeCacheSTG]
directory=/home/ec2-user/MMSys-ws/Servers
command=/home/ec2-user/anaconda3/envs/lbk-tornado/bin/python stg_trade_cache.py
autostart=true
autorestart=unexpected
startretries=5
startsecs=1
priority=999
stderr_logfile=/dev/null
stdout_logfile=/dev/null
stopasgroup=true
killasgroup=true
user = root
redirect_stderr = true
asgfjfghfc_maxbytes = 20M
stdout_logfile_backups = 20
